# crontab
# Minuto  Hora  Dia  Mes  Dia_semana  Comando
35 17 * * 1-5 /usr/bin/python3 script 
